# service/upstart/amasel.conf
# amasel - Sync daemon for Amasel service
#

description     "Amasel Sync daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022


exec /srv/amasel/worker/run_amasel_syncer.sh